# systemd/pixeldock32.service
[Unit]
Description=PixelDock32 FastAPI LED Controller
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/PixelDock32
EnvironmentFile=/home/pi/PixelDock32/.env
ExecStart=/home/pi/PixelDock32/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
